# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:jenkins]
command=bash -c "/etc/init.d/jenkins start"
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
